# dendrite.cron
#!/bin/sh
#
# This file needs to be installed in /etc/cron.hourly
# in order to kick the dendrite to check itself every hour 
# 

/usr/sbin/dendrite
